# containers/rpi/files/image-inference/systemd/squirrel-away-inference-image.path
[Unit]
Description=Watch squirrel-away image capture directory for new images.

[Path]
DirectoryNotEmpty=/var/lib/squirrel-away/inference/input/
MakeDirectory=true

[Install]
WantedBy=default.target
